# deploy/tf-dashboard.service
[Unit]
Description=TF Dashboard Web Service
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/tfamily-fans-monitor
Environment="PYTHONUNBUFFERED=1"
Environment="TF_DB_FILE=/home/ubuntu/tfamily-fans-monitor/data/tf_dashboard.sqlite3"
Environment="TF_DASHBOARD_HOST=127.0.0.1"
Environment="TF_DASHBOARD_PORT=8000"
ExecStart=/home/ubuntu/tfamily-fans-monitor/.venv/bin/python /home/ubuntu/tfamily-fans-monitor/dashboard_server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
