# deploy/tf-spider.service
[Unit]
Description=TF Spider Collector
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu/tf-dashboard
Environment="PYTHONUNBUFFERED=1"
Environment="TF_AUTH_TOKEN=replace-with-your-token"
Environment="TF_DB_FILE=/home/ubuntu/tf-dashboard/data/tf_dashboard.sqlite3"
Environment="TF_POLL_INTERVAL_SECONDS=150"
ExecStart=/home/ubuntu/tf-dashboard/.venv/bin/python /home/ubuntu/tf-dashboard/spider.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
